# deploy/systemd/rob-portal-dev.service
[Unit]
Description=Rob Portal Dev Service
After=network.target

[Service]
Type=simple
User=deployuser
Group=deployuser
WorkingDirectory=/opt/rob-portal/app/portal
EnvironmentFile=/opt/rob-portal/app/.env
ExecStart=/opt/rob-portal/app/.venv/bin/gunicorn rob_portal.wsgi:application --bind 127.0.0.1:8090
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
